# files/etc/uci-defaults/04_gl_upgrade
uci set upgrade.general.url=www
uci commit upgrade
